# supervisord.conf
[program:daemonname]
command = python /home/dotcloud/current/phew.py